# launchd/com.local.litellm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.local.litellm</string>
  <key>RunAtLoad</key>      <true/>
  <key>KeepAlive</key>      <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>      <string>@@REPO_ROOT@@/.venvs/litellm/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>      <string>@@HOME@@</string>
    <key>USER</key>      <string>@@USER@@</string>
    <key>TMPDIR</key>    <string>@@TMPDIR@@</string>
    <key>PYTHONPATH</key><string>@@REPO_ROOT@@</string>
  </dict>
  <!-- See com.local.ollama.plist: logs are redirected by the child shell
       because launchd cannot open files under ~/Documents (TCC). -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd "@@REPO_ROOT@@" &amp;&amp; exec "@@REPO_ROOT@@/.venvs/litellm/bin/python" "@@REPO_ROOT@@/scripts/run_litellm.py" --config "@@REPO_ROOT@@/config/litellm-config.rendered.yaml" --host 127.0.0.1 --port 4000 &gt;&gt; "@@REPO_ROOT@@/.logs/litellm.out.log" 2&gt;&gt; "@@REPO_ROOT@@/.logs/litellm.err.log"</string>
  </array>
</dict>
</plist>
