# launchd/com.local.ollama-warm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>           <string>com.local.ollama-warm</string>
  <key>RunAtLoad</key>       <true/>
  <key>KeepAlive</key>       <false/>
  <key>WorkingDirectory</key><string>@@REPO_ROOT@@</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>          <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OLLAMA_HOST</key>   <string>127.0.0.1:@@OLLAMA_PORT@@</string>
    <key>OLLAMA_TAG</key>    <string>@@OLLAMA_TAG@@</string>
    <key>WARM_WAIT_S</key>   <string>180</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@REPO_ROOT@@/scripts/warm-ollama.sh</string>
  </array>
  <key>StandardOutPath</key><string>@@REPO_ROOT@@/.logs/ollama-warm.out.log</string>
  <key>StandardErrorPath</key><string>@@REPO_ROOT@@/.logs/ollama-warm.err.log</string>
</dict>
</plist>
