# launchd/com.local.ollama.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.local.ollama</string>
  <key>RunAtLoad</key>      <true/>
  <key>KeepAlive</key>      <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>                   <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>                   <string>@@HOME@@</string>
    <key>USER</key>                   <string>@@USER@@</string>
    <key>TMPDIR</key>                 <string>@@TMPDIR@@</string>
    <key>OLLAMA_HOST</key>            <string>127.0.0.1:11434</string>
    <key>OLLAMA_KV_CACHE_TYPE</key>   <string>@@KV_CACHE_TYPE@@</string>
    <key>OLLAMA_FLASH_ATTENTION</key> <string>1</string>
    <key>OLLAMA_NUM_PARALLEL</key>    <string>@@OLLAMA_NUM_PARALLEL@@</string>
    <key>OLLAMA_KEEP_ALIVE</key>      <string>@@OLLAMA_KEEP_ALIVE@@</string>
  </dict>
  <!-- The repo lives under ~/Documents, a TCC-protected path that launchd
       cannot open log files in. So we omit StandardOutPath/StandardErrorPath
       (and WorkingDirectory) and let the child shell, which runs with the
       user's TCC context, cd into the repo and redirect into .logs itself. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd "@@REPO_ROOT@@" &amp;&amp; exec /opt/homebrew/bin/ollama serve &gt;&gt; "@@REPO_ROOT@@/.logs/ollama.out.log" 2&gt;&gt; "@@REPO_ROOT@@/.logs/ollama.err.log"</string>
  </array>
</dict>
</plist>
