# launchd/com.local.mlx.rendered.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>          <string>com.local.mlx</string>
  <key>RunAtLoad</key>      <true/>
  <key>KeepAlive</key>      <true/>
  <key>ProcessType</key>    <string>Interactive</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>   <string>/Users/martinfr/Documents/GitHub/MacM4LocalAgent/.venvs/mlx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>   <string>/Users/martinfr</string>
    <key>USER</key>   <string>martinfr</string>
    <key>TMPDIR</key> <string>/var/folders/6c/pyk36nb12vn1t4ynjrq_hkth0000gn/T/</string>
  </dict>
  <!-- See com.local.ollama.plist: logs are redirected by the child shell
       because launchd cannot open files under ~/Documents (TCC). -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd "/Users/martinfr/Documents/GitHub/MacM4LocalAgent" &amp;&amp; exec "/Users/martinfr/Documents/GitHub/MacM4LocalAgent/.venvs/mlx/bin/python" -m mlx_lm server --model "/Users/martinfr/Documents/GitHub/MacM4LocalAgent/models/mlx-community_Qwen3-Coder-Next-4bit" --host 127.0.0.1 --port 8081 &gt;&gt; "/Users/martinfr/Documents/GitHub/MacM4LocalAgent/.logs/mlx.out.log" 2&gt;&gt; "/Users/martinfr/Documents/GitHub/MacM4LocalAgent/.logs/mlx.err.log"</string>
  </array>
</dict>
</plist>
