# launchd/com.local.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>             <string>com.local.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@REPO_ROOT@@/scripts/95-watchdog.sh</string>
  </array>

  <!-- Fire every 5 minutes to clean up orphaned LiteLLM tasks / Ollama generations -->
  <key>StartInterval</key>     <integer>300</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>   <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>   <string>@@REPO_ROOT@@/.logs/watchdog.log</string>
  <key>StandardErrorPath</key> <string>@@REPO_ROOT@@/.logs/watchdog.log</string>
  <key>RunAtLoad</key>         <false/>
</dict>
</plist>
